# docker/supervisord.conf
[supervisord]
nodaemon=true
user=www-data

[program:php-fpm]
command=docker-php-entrypoint php-fpm

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
